# mastodon_supervisor.ini
[group:mastodon]
programs=web,sidekiq_default,sidekiq_scheduler,sidekiq_push,sidekiq_pull,streaming

[program:web]
command=/Users/foo/.gem/ruby/3.2.2/bin/bundle exec puma -C config/puma.rb
directory=/Users/foo/Sites/mastodon
stdout_logfile=/Users/foo/Sites/mastodon/log/puma.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
redirect_stderr=true
environment=RAILS_ENV=production,PORT=3000,RAILS_LOG_LEVEL=warn,MAX_THREADS=10,WEB_CONCURRENCY=4,PGGSSENCMODE=disable,PATH="/opt/homebrew/bin:%(ENV_PATH)s"
stopasgroup=true

[program:sidekiq_default]
command=/Users/foo/.gem/ruby/3.2.2/bin/bundle exec sidekiq -c 30 -q default -q mailers -q push -q pull -q ingress
directory=/Users/foo/Sites/mastodon
stdout_logfile=/Users/foo/Sites/mastodon/log/sidekiq.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
redirect_stderr=true
environment=RAILS_ENV=production,DB_POOL=30,PGGSSENCMODE=disable
stopasgroup=true

[program:sidekiq_scheduler]
command=/Users/foo/.gem/ruby/3.2.2/bin/bundle exec sidekiq -c 5 -q scheduler -q default -q mailers -q push -q pull -q ingress
directory=/Users/foo/Sites/mastodon
stdout_logfile=/Users/foo/Sites/mastodon/log/sidekiq.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
redirect_stderr=true
environment=RAILS_ENV=production,DB_POOL=5,PGGSSENCMODE=disable
stopasgroup=true

[program:sidekiq_pull]
command=/Users/foo/.gem/ruby/3.2.2/bin/bundle exec sidekiq -c 30 -q pull -q ingress -q push -q default -q mailers
directory=/Users/foo/Sites/mastodon
stdout_logfile=/Users/foo/Sites/mastodon/log/sidekiq.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
redirect_stderr=true
environment=RAILS_ENV=production,DB_POOL=30,PGGSSENCMODE=disable
stopasgroup=true

[program:sidekiq_push]
command=/Users/foo/.gem/ruby/3.2.2/bin/bundle exec sidekiq -c 30 -q push -q pull -q ingress -q default -q mailers
directory=/Users/foo/Sites/mastodon
stdout_logfile=/Users/foo/Sites/mastodon/log/sidekiq.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
redirect_stderr=true
environment=RAILS_ENV=production,DB_POOL=30,PGGSSENCMODE=disable
stopasgroup=true

[program:streaming]
command=/Users/foo/.nodenv/shims/npm run start
directory=/Users/foo/Sites/mastodon
stdout_logfile=/Users/foo/Sites/mastodon/log/streaming.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
redirect_stderr=true
environment=NODE_ENV=production,PORT=4000
stopasgroup=true